# crontab
* * * * * php /var/www/artisan schedule:run >> /root/laravel-schedule.log 2>&1
