# worker.ini
[supervisord]
nodaemon=true

[program:nginx]
command = /usr/sbin/nginx
user = root
autostart = true
autorestart = true

[program:php-fpm]
command = /usr/sbin/php-fpm7 -F -R
user = root
autostart = true
autorestart = true

[program:crond]
command=/usr/sbin/crond -f
autostart = true
autorestart = true

[program:laravel-worker]
command=php /var/www/artisan queue:work --sleep=10 --quiet --timeout=60 --tries=3
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=3
stdout_logfile=/root/laravel-worker.log
