# worker.ini

[program:nginx]
command = /usr/sbin/nginx
user = root
autostart = true

[program:php-fpm]
command = /usr/bin/php-fpm -F
user = root
autostart = true

[program:laravel-worker]
command=php /var/www/artisan queue:work --sleep=10 --quiet --timeout=60 --tries=3
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=3
stdout_logfile=/root/laravel-worker.log
